# m2o_live/supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:live_nginx]
command=/usr/local/Hoge/sbin/nginx

[program:ts_manager]
command=/usr/local/Hoge/bin/ts_manager_run 2> &1

[program:live_db]
command=cd /m2odata/server/pyrtmp && ./manage  db upgrade

[program:pyrtmp]
command=sh -c "PYRTMP_CONFIG=/m2odata/server/pyrtmp/pyrtmp.conf /m2odata/server/pyrtmp/manage runserver"
autostart=true
stopsignal=INT
stopwaitsecs=30
log_stdout=true
log_stderr=true
logfile=/var/log/supervisor/pyrtmp.log
logfile_maxbytes=1MB
logfile_backups=10
