# m2o_vod_php/supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:vod_nginx]
command=/etc/init.d/vod_nginx start

[program:php-fpm]
command=/etc/init.d/php-fpm start
